# raspberrypios_battery/update_battery.service
[Unit]
Description=Update virtual battery level

[Service]
Type=oneshot
ExecStart=/usr/bin/update_battery.sh
